# scripts/comms_synth.service
[Unit]
Description=Comms365 Synth
After=network.target

[Service]
Type=simple
Environment=STAGE=production
WorkingDirectory=/var/www/synth/current
ExecStart=/usr/bin/python3 synth OnCommsDevicePilot demo_dashboard_sensors
Restart=on-failure
User=pervasive
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cityfibre_synth

[Install]
WantedBy=multi-user.target
